# .claude/launchd/com.wally.dashboard.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.wally.dashboard</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/josecampos/Documents/wally-trader/shared/wally_core/.venv/bin/python</string>
        <string>-m</string>
        <string>wally_core.dashboard_server</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>WorkingDirectory</key>
    <string>/Users/josecampos/Documents/wally-trader</string>
    <key>StandardOutPath</key>
    <string>/Users/josecampos/Documents/wally-trader/logs/dashboard.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/josecampos/Documents/wally-trader/logs/dashboard.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONPATH</key>
        <string>/Users/josecampos/Documents/wally-trader/shared/wally_core/src</string>
    </dict>
</dict>
</plist>
