# .claude/launchd/com.wally.funding-monitor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.wally.funding-monitor</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/env</string>
        <string>python3</string>
        <string>/Users/josecampos/Documents/wally-trader/.claude/scripts/funding_alerts.py</string>
        <string>--once</string>
        <string>--notify</string>
    </array>
    <key>StartInterval</key>
    <integer>3600</integer>
    <key>StandardOutPath</key>
    <string>/Users/josecampos/Documents/wally-trader/logs/funding-monitor.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/josecampos/Documents/wally-trader/logs/funding-monitor.log</string>
    <key>WorkingDirectory</key>
    <string>/Users/josecampos/Documents/wally-trader</string>
</dict>
</plist>
