# .claude/launchd/com.wally.post-mortem-watcher.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.wally.post-mortem-watcher</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/josecampos/Documents/wally-trader/shared/wally_core/.venv/bin/python</string>
        <string>/Users/josecampos/Documents/wally-trader/.claude/scripts/learning/post_mortem_trigger.py</string>
        <string>--profile</string>
        <string>bitunix</string>
        <string>--watch</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/josecampos/Documents/wally-trader</string>

    <key>KeepAlive</key>
    <true/>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/Users/josecampos/Library/Logs/wally-trader/post-mortem-watcher.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/josecampos/Library/Logs/wally-trader/post-mortem-watcher.err</string>

    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
